# ansible/roles/restic/templates/restic-backup.timer.j2
[Unit]
Description=Daily restic backup for {{ repo.name }}

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
